# fakeSMTP@.service
[Unit]
Description=FakeSMTP Script

[Service]
ExecStart=/usr/local/bin/nitmo_fakesmtp.sh
StandardInput=socket
